# python-broker/deploy/vcw-flow-broker.service
[Unit]
Description=VCW Google Flow Python Broker (invisible_playwright x86_64 via FEX-Emu)
After=network.target xvfb.service
Wants=xvfb.service

[Service]
Type=simple
User=truonghoanghuy
Group=truonghoanghuy
WorkingDirectory=/opt/vcw/app/python-broker
# Pull BROKER_AUTH_TOKEN + any overrides from this file.
# Server's existing .env stays untouched.
EnvironmentFile=-/opt/vcw/app/python-broker/broker.env
# Headful Firefox renders on the external Xvfb display (xvfb.service, :99).
Environment=DISPLAY=:99
# PROVEN engine: invisible_playwright (x86_64 stealth Firefox) running under
# FEX-Emu on this aarch64 box. Camoufox was tried and FAILS reCAPTCHA on Flow
# (grecaptcha never initializes) — DO NOT switch to camoufox here. HOME is set
# explicitly so the FEX rootfs + Firefox cache resolve under ~truonghoanghuy.
Environment=BROKER_BROWSER_ENGINE=invisible
Environment=HOME=/home/truonghoanghuy
ExecStart=/opt/vcw/app/python-broker/venv-x86/bin/python3 -m broker
Restart=on-failure
RestartSec=5
StandardOutput=append:/opt/vcw/logs/broker.log
StandardError=append:/opt/vcw/logs/broker-error.log

# Resource guard — Firefox under FEX (parent + content procs) can balloon.
# Box has 23G RAM + 4G swap, so headroom above the old 2G ceiling is safe.
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=multi-user.target
